# Scheduled Antivirus Scans Scripts/crontab
0 5 * * 1 echo "/home/pi/" "/home/pi/Documents/Programming/automation-sysadmin-scripts/Scheduled Antivirus Scans Scripts/scan_directory.sh"
